# cems-service-common/pack/installZip/bin/CEMS-SERVICE-SCAN
#!/bin/sh
#
# "$Id: cups.sh,v 1.10 2014/11/12 9:15:16 MJD $"
#
#   Startup/shutdown script for the Common UNIX Printing System (CEMS-SERVICE-SCAN).
#
#   Linux chkconfig stuff:
#
#   chkconfig: 2345 98 50
#   description: Startup/shutdown script for the Common UNIX \
#                Printing System (CEMS-SERVICE-SCAN).
#
#   Copyright 2010-2015 by vrv co.,ltd, all rights reserved.
#
#   These coded instructions, statements, and computer programs are the
#   property of VRV Software Products Distribution and use rights are outlined in the file
#   "LICENSE.txt" which should have been included with this file.  If this
#   file is missing or damaged
#
#       WWW: http://www.vrv.com.cn
#
# heavily edited so that it's more like other scripts in init.d on Red Hat
# Linux
#
### BEGIN INIT INFO
# Provides: CEMS-SERVICE-SCAN
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: The CEMS-SERVICE-SCAN scheduler
# Description: The CEMS-SERVICE-SCAN scheduler
### END INIT INFO

# Source function library.
. /etc/rc.d/init.d/functions

ServiceName=CEMS-SERVICE-SCAN
DAEMON_DIR="/usr/local/service/CEMS-SERVICE-SCAN/bin"
proc_name="cemsscan"
DAEMON="$DAEMON_DIR/$proc_name"

RETVAL=0

function is_runing()
{
    pid=$(pgrep $proc_name)
    if [ "$pid" != "" ] 
    then
        return 0
    fi
    return 1
}

start () {
	# start daemon
    [ -d "$DAEMON_DIR" ] && cd $DAEMON_DIR
	echo -n $"Starting $ServiceName..."
	$DAEMON
    RETVAL=$?
    sleep 1
    is_runing
    if [ "$?" = "0" ]
    then
        echo -e "\e[1;32mStart success!\e[0m"
    else
        echo -e "\e[1;31mStart failed!\e[0m"
    fi
}

stop () {
	# stop daemon
    is_runing
    if [ "$?" = "0" ]
    then
        pgrep $proc_name | xargs -t -I{} kill -9 {} 1>&2 2>/dev/null
    else
        echo -e "\e[1;32m$ServiceName has stoped!\e[0m"
        return 0
    fi

    echo -e "\e[1;32m$ServiceName stop success!\e[0m"
    return 0
}

restart() {
	stop
	start
}

status() {
		local result
		echo "check status of ${ServiceName}..."

		result=$(pgrep $proc_name)

		if [ "$result" != "" ]
        then
			echo -e "\e[1;32m${ServiceName} is up\e[0m"

			RETVAL=0

		else
			echo -e "\e[1;31m${ServiceName} is down\e[0m"

			RETVAL=1

		fi

		echo "check status of ${ServiceName} ... done"

                RETVAL=$?
		
}

case $1 in
	start)
		start
	;;
	stop)
		stop
	;;
	restart)
		restart
	;;	
	status)
		status
	;;
	*)

	echo $"Usage: $ServiceName {start|stop|restart|status}" 
	exit 1
esac

exit $RETVAL
